# server/cookbooks/iperf3/files/etc/systemd/system/iperf3.service
[Unit]
Description=iperf3
After=network.target

[Service]
User=nobody
ExecStart=/usr/bin/iperf3 -s

[Install]
WantedBy=multi-user.target
