# server/cookbooks/systemd-nspawn/templates/etc/systemd/system/systemd-nspawn@.service.d/bridge.conf
[Service]
ExecStart=
ExecStart=/usr/bin/systemd-nspawn --quiet --keep-unit --boot --link-journal=try-guest --network-bridge=<%= node[:systemd_nspawn][:bridge] %> -U --settings=override --machine=%i
